# files/systemd/ppstime-healthcheck.service
[Unit]
Description=PPSPi passive timing health check
After=chrony.service gpsd.service

[Service]
Type=oneshot
ExecStart=/usr/lib/ppstime/ppstime-healthcheck
Nice=10