# files/systemd/ppstime-healthcheck.timer
[Unit]
Description=Run the PPSPi health check without restarting timing services

[Timer]
OnBootSec=5min
OnUnitActiveSec=2min
RandomizedDelaySec=15s
AccuracySec=15s
Unit=ppstime-healthcheck.service

[Install]
WantedBy=timers.target
