# files/systemd/ppstime-healthcheck.service
[Unit]
Description=PPSPi passive timing health check
After=chrony.service gpsd.service

[Service]
Type=oneshot
ExecStart=/usr/lib/ppstime/ppstime-healthcheck
TimeoutStartSec=45s
Nice=10
UMask=0027
RuntimeDirectory=ppstime
RuntimeDirectoryMode=0755
RuntimeDirectoryPreserve=yes
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectKernelLogs=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
