# files/systemd/ppstime-healthcheck.timer
[Unit]
Description=Run the PPSPi health check without restarting timing services

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
RandomizedDelaySec=1min
Unit=ppstime-healthcheck.service

[Install]
WantedBy=timers.target