# files/systemd/gpsd.service.d/ppstime.conf
[Unit]
Wants=chrony.service
After=chrony.service systemd-udev-settle.service

[Service]
ExecStartPre=/usr/lib/ppstime/ppstime-wait-devices
Restart=on-failure
RestartSec=5s