# files/systemd/chrony.service.d/ppstime.conf
[Unit]
After=ppstime-rtc-restore.service
Before=gpsd.service